# .snapshots/auto-2026-06-17_01-35-59/config/battery-charge-hysteresis/battery-charge-hysteresis.timer
[Unit]
Description=Run battery charge hysteresis threshold check

[Timer]
OnBootSec=30s
OnUnitActiveSec=20min
AccuracySec=5min
Persistent=true
Unit=battery-charge-hysteresis.service

[Install]
WantedBy=timers.target
