# config/battery-charge-probe/battery-charge-probe.service
[Unit]
Description=Battery charge threshold probe
ConditionPathExists=/sys/class/power_supply/BAT0/charge_control_end_threshold

[Service]
Type=oneshot
Environment=BATTERY_PATH=/sys/class/power_supply/BAT0
Environment=LOW_THRESHOLD=40
Environment=HIGH_THRESHOLD=60
ExecStart=/usr/local/bin/battery-charge-probe
